# PiGi.service
[Unit]
Description=PiGi service

[Service]
Type=simple
ExecStart=/usr/bin/python /home/nymea/PiGi/PiGi.py

[Install]
WantedBy=multi-user.target
